# aiot/supervisord-mysql.ini
[program:chownmysql]
command=/etc/initMysql.sh

[program:mysql]
command=/usr/mysql/mysql-5.7.22-linux-glibc2.12-x86_64/support-files/mysql.server start
